# ttyPS1.conf
start on runlevel [23]
stop on runlevel [!23]
respawn
exec /sbin/getty -8 115200 ttyPS1
